# prog/sys/busybox/cfg/sysinit-latest/etc/inittab
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == sysinit, respawn, askfirst, wait, once, ctrlaltdel or shutdown
# process   == program to run

::sysinit:/etc/init.d/rcS

tty1::respawn:/bin/busybox getty 38400 tty1 linux
tty2::respawn:/bin/busybox getty 38400 tty2 linux
tty3::respawn:/bin/busybox getty 38400 tty3 linux

::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK

