# prog/sys/logrotate/cfg/all-latest/etc/cron.d/logrotate
*/30 * * * * root /sbin/logrotate -l syslog /etc/logrotate.conf
